# scripts/com.ryanileo.slack-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ryanileo.slack-listener</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ahnbingbing/code/rianileo-agent/.venv/bin/python</string>
        <string>/Users/ahnbingbing/code/rianileo-agent/scripts/slack_listener.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ahnbingbing/code/rianileo-agent</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/ahnbingbing/code/rianileo-agent/data/logs/slack_listener.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ahnbingbing/code/rianileo-agent/data/logs/slack_listener.err.log</string>
</dict>
</plist>
