# launchd/com.rianileo.giri-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Giri review-learning — fires every 3 days (FULL AUTO, PD 2026-06-21, temporary until
  the pipeline stabilizes). Runs headless Claude Code locally so it can read the local
  Slack/DB/working-tree and edit Giri. Edits left uncommitted + summarized to the board.

  INSTALL (PD runs this — it authorizes the unsupervised auto-editor):
    cp launchd/com.rianileo.giri-weekly.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.rianileo.giri-weekly.plist
  DISABLE:
    launchctl bootout gui/$(id -u)/com.rianileo.giri-weekly
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rianileo.giri-weekly</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/ahnbingbing/code/rianileo-agent/scripts/giri_weekly_update.sh</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/ahnbingbing/code/rianileo-agent</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>TZ</key><string>Asia/Seoul</string>
      <key>PATH</key><string>/Users/ahnbingbing/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <!-- every 3 days (259200s); does NOT fire on load (RunAtLoad omitted) -->
    <key>StartInterval</key><integer>259200</integer>
    <key>StandardOutPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/giri_weekly.out.log</string>
    <key>StandardErrorPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/giri_weekly.err.log</string>
    <key>ProcessType</key><string>Background</string>
  </dict>
</plist>
