# launchd/com.rianileo.writer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily Writer agent run @ 18:00 KST.

  Install:
    cp launchd/com.rianileo.writer.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.rianileo.writer.plist

  Adjust paths below before installing:
    /Users/ligi/code/rianileo-agent  -> your actual repo path
    /Users/ligi/.pyenv/shims/python  -> your Python (or .venv/bin/python)
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rianileo.writer</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/ligi/code/rianileo-agent/.venv/bin/python</string>
      <string>-m</string>
      <string>agents.writer</string>
    </array>

    <key>WorkingDirectory</key><string>/Users/ligi/code/rianileo-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>TZ</key><string>Asia/Seoul</string>
      <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key><string>/Users/ligi/code/rianileo-agent/data/logs/writer.out.log</string>
    <key>StandardErrorPath</key><string>/Users/ligi/code/rianileo-agent/data/logs/writer.err.log</string>

    <key>RunAtLoad</key><false/>
    <key>KeepAlive</key><false/>
  </dict>
</plist>
