# launchd/com.rianileo.gmp-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- grandmompapa proactive nudge (morning, 9:00 KST). RunAtLoad off so loading
     never posts to the channel. Install:
       cp launchd/com.rianileo.gmp-morning.plist ~/Library/LaunchAgents/
       launchctl unload ~/Library/LaunchAgents/com.rianileo.gmp-morning.plist 2>/dev/null
       launchctl load   ~/Library/LaunchAgents/com.rianileo.gmp-morning.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rianileo.gmp-morning</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/ahnbingbing/code/rianileo-agent/.venv/bin/python</string>
      <string>scripts/grandmompapa_nudge.py</string>
      <string>morning</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/ahnbingbing/code/rianileo-agent</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>TZ</key><string>Asia/Seoul</string>
      <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
      <key>LOG_LEVEL</key><string>INFO</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <key>StandardOutPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/gmp_nudge.out.log</string>
    <key>StandardErrorPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/gmp_nudge.err.log</string>
    <key>ProcessType</key><string>Background</string>
  </dict>
</plist>
