# launchd/com.rianileo.gcs-backfill-oneshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ONE-SHOT: backfill the 735 offloaded originals into GCS (PD 2026-06-28).

  Why launchd and not a plain shell run: osxphotos needs Photos-library TCC
  authorization, which is granted to the .venv python ONLY when launchd is the
  responsible parent (the dawn icloud-sync job has it; an ad-hoc nohup from a
  terminal/Claude session does NOT → "could not get authorization").

  Run once:
      cp launchd/com.rianileo.gcs-backfill-oneshot.plist ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.rianileo.gcs-backfill-oneshot.plist
  Tail:
      tail -f data/logs/gcs_backfill_769.log
  When DONE, remove it (it's not a recurring job):
      launchctl unload ~/Library/LaunchAgents/com.rianileo.gcs-backfill-oneshot.plist
      rm ~/Library/LaunchAgents/com.rianileo.gcs-backfill-oneshot.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rianileo.gcs-backfill-oneshot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/ahnbingbing/code/rianileo-agent/.venv/bin/python</string>
      <string>-m</string>
      <string>scripts.gcs_backfill</string>
      <string>--batch</string>
      <string>60</string>
    </array>

    <key>WorkingDirectory</key><string>/Users/ahnbingbing/code/rianileo-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>TZ</key><string>Asia/Seoul</string>
      <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
      <key>LOG_LEVEL</key><string>INFO</string>
    </dict>

    <key>StandardOutPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/gcs_backfill_769.log</string>
    <key>StandardErrorPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/gcs_backfill_769.log</string>

    <key>RunAtLoad</key><true/>
    <key>ProcessType</key><string>Background</string>
    <key>Nice</key><integer>10</integer>
  </dict>
</plist>
