# launchd/com.rianileo.bandit-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  av-vs-rf bandit — WEEKLY Slack report (PD 2026-06-07, first_month_plan §4).

  Posts "주간 av vs rf 누적 reward·표본수·현재 우세확률" to the workroom every
  Monday 10:00 KST (Weekday 1). Also collects fresh metrics first.

  Install:
      cp launchd/com.rianileo.bandit-report.plist ~/Library/LaunchAgents/
      launchctl unload ~/Library/LaunchAgents/com.rianileo.bandit-report.plist 2>/dev/null
      launchctl load   ~/Library/LaunchAgents/com.rianileo.bandit-report.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rianileo.bandit-report</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/ahnbingbing/code/rianileo-agent/.venv/bin/python</string>
      <string>-m</string>
      <string>agents.bandit</string>
      <string>--collect</string>
      <string>--slack</string>
    </array>

    <key>WorkingDirectory</key><string>/Users/ahnbingbing/code/rianileo-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>TZ</key><string>Asia/Seoul</string>
      <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
      <key>LOG_LEVEL</key><string>INFO</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/bandit.out.log</string>
    <key>StandardErrorPath</key><string>/Users/ahnbingbing/code/rianileo-agent/data/logs/bandit.err.log</string>

    <key>RunAtLoad</key><false/>
    <key>ProcessType</key><string>Background</string>
    <key>Nice</key><integer>10</integer>
  </dict>
</plist>
